# email-validator-pro/systemd/email-validator.service
# /etc/systemd/system/email-validator.service
[Unit]
Description=EmailValidator Pro Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/email-validator-pro
ExecStart=/home/ubuntu/email-validator-pro/venv/bin/python app.py
Restart=always
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
